# vps-proxy/systemd.service
[Unit]
Description=Supabase RU proxy (Node)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/app/vps-proxy
Environment=PORT=8787
Environment=SUPABASE_URL=https://jvibhsjnspvucjwvhfht.supabase.co
ExecStart=/usr/bin/node /opt/app/vps-proxy/server.js
Restart=always
RestartSec=2
User=www-data

[Install]
WantedBy=multi-user.target